# homebox/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Homebox
# Runs the NGINX daemon
# ==============================================================================
# Wait for Homebox to become available
bashio::net.wait_for 7745 localhost 90

bashio::log.info "Starting NGINX..."

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi
dns_host=$(bashio::dns.host)
sed -i "s/%%dns_host%%/${dns_host}/g" /etc/nginx/includes/resolver.conf

exec nginx
