# homebox/rootfs/etc/services.d/homebox/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Homebox
# Runs the Homebox daemon
# ==============================================================================

bashio::log.info "Starting Homebox..."


exec /app/api
